# numLockOnTty.service
[Unit]
Description=numlockOnTty

[Service]
ExecStart=/usr/bin/numlockOnTty

[Install]
WantedBy=multi-user.target

